# .config/systemd/user/night-light.timer
[Unit]
Description=Run night-light every 60 min

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
